# brew/templates/update-brewfile.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
   "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
    <string>@LABEL@</string>
  <key>ProgramArguments</key>
    <array>
      <string>@SCRIPT_PATH@</string>
    </array>
  <key>StartCalendarInterval</key>
    <dict>
      <key>Weekday</key><integer>@WEEKDAY@</integer>   <!-- Mon=2 -->
      <key>Hour</key><integer>@HOUR@</integer>
      <key>Minute</key><integer>@MINUTE@</integer>
    </dict>
  <key>StandardOutPath</key>
    <string>@LOG_OUT@</string>
  <key>StandardErrorPath</key>
    <string>@LOG_ERR@</string>
</dict>
</plist>
